# deploy/reverb.supervisor.conf
; Supervisor configuration for Progressia Reverb WebSocket server
;
; Place this file at /etc/supervisor/conf.d/progressia-reverb.conf
; (or another path included by your supervisord.conf), then run:
;
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start progressia-reverb
;
; Adjust `user`, `command`, and `directory` to match your deploy setup.
; In production, front Reverb with nginx on port 443 (path /app) so
; the browser can use wss://your-domain.example/app.

[program:progressia-reverb]
process_name=%(program_name)s
command=php /var/www/progressia/artisan reverb:start --port=8080 --hostname=0.0.0.0
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/progressia/storage/logs/reverb.log
stopwaitsecs=10
